# crabby-secrets.service
[Unit]
Description=Crabby Secrets API
After=network.target

[Service]
User=crabby
WorkingDirectory=/home/crabby/secrets-api
ExecStart=/home/crabby/secrets-api/venv/bin/gunicorn --bind 127.0.0.1:5002 --workers 1 app:app
Restart=always
RestartSec=5
Environment=KEYS_FILE=/home/crabby/.secrets/keys.json
Environment=TOKEN_FILE=/home/crabby/.secrets/api-token.txt

[Install]
WantedBy=multi-user.target
